# deploy/chatter-frontend.service
# Template for manual installation. `deploy/install.sh` can generate a customised unit automatically.

[Unit]
Description=Chatter BBS web control deck
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/chatter-web
ExecStart=/usr/bin/node /opt/chatter-web/server.js
Restart=on-failure
Environment=PORT=8081
Environment=CHATTER_BBS_HOST=bbs.chatter.example
Environment=CHATTER_BBS_PROTOCOL=telnet
Environment=CHATTER_BBS_PORT=23
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
